# cross-toolchain/yocto-sdk/sysroots/core2-32-sbr-linux/nao/etc/systemd/user/crash-report-uploader.service
[Unit]
Description=Uploader for NAOqi minidump
After=network.target network-online.target local-fs.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStartPre=/usr/libexec/aldebaran/clean-dump-dir.sh
ExecStart=/opt/aldebaran/bin/crash-report-upload
