# cross-toolchain/yocto-sdk/sysroots/core2-32-sbr-linux/etc/systemd/system/data.target
[Unit]
After=data-skeleton.service
Requires=data-skeleton.service data-recover.service
Wants=data-repair.service
