# cross-toolchain/yocto-sdk/sysroots/core2-32-sbr-linux/etc/systemd/system/data-repair.service
[Unit]
Description=verify and repair user partition
RefuseManualStart=true
DefaultDependencies=false
After=local-fs-pre.target tmp.mount systemd-udevd.service dev-disk-by\x2duuid-66666666\x2d6666\x2d6666\x2d6666\x2d666666666666.device
Requires=tmp.mount
Wants=dev-disk-by\x2duuid-66666666\x2d6666\x2d6666\x2d6666\x2d666666666666.device
OnFailure=notification@205.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/aldebaran/checksd.sh --repair
TimeoutSec=3600
