# cross-toolchain/yocto-sdk/sysroots/core2-32-sbr-linux/etc/systemd/system/hostname.service
[Unit]
Description=Hostname
After=local-fs.target
Before=network.target
[Service]
Type=oneshot
ExecStart=/bin/hostname -F /etc/hostname
