# cross-toolchain/yocto-sdk/sysroots/core2-32-sbr-linux/etc/systemd/system/media-internal.target
[Unit]
After=media-internal.mount media-internal-skeleton.service
Requires=media-internal.mount media-internal-skeleton.service
