# cross-toolchain/yocto-sdk/sysroots/core2-32-sbr-linux/etc/systemd/system/user@1001.service.d/env.conf
[Service]
# See spec for XDG_CONFIG_DIRS: http://standards.freedesktop.org/basedir-spec/basedir-spec-latest.html
Environment=XDG_CONFIG_DIRS=/opt/aldebaran/etc:/nao/etc:/etc/xdg
